# roles/glauth/templates/glauth-proxy.socket.j2
[Socket]
ListenStream={{ glauth_port | default(389) }}

[Install]
WantedBy=sockets.target
